# crontab
0,5,10,15,20,25,30,35,41,45,50,55 * * * * DEPLOY_PATH/scripts/notify-users.sh >> /var/log/cron.log 2>&1
0 0 * * 3 DEPLOY_PATH/scripts/clear-log.sh >> /var/log/cron.log 2>&1
